# packaging/systemd/deepseek-runtime.service
[Unit]
Description=DeepSeekCode HTTP runtime
After=network.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
ExecStart=/usr/bin/env __DEEPSEEK_BIN__ serve --http --addr __HTTP_ADDR__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
