# packaging/systemd/deepseek-shell-supervisor.service
[Unit]
Description=DeepSeekCode shell supervisor protocol bridge
# Exposes workspace-local shell status/show/start/control; native PTY sessions are not implemented yet.
After=network.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
ExecStart=/usr/bin/env __DEEPSEEK_BIN__ agents shell-supervisor --json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
